# sys-apps/my-systemd-units/files/services-server/nfs4_client.target
[Unit]
Description=Dependencies for NFSv4 mounts
Wants=rpcbind.service
Wants=nfs-utils.service

[Install]
WantedBy=remote-fs.target
