# net-p2p/jackett/files/jackett.socket
[Unit]
Description=Jackett search service
Conflicts=jackett.service

[Socket]
ListenStream=9117
Accept=yes

[Install]
WantedBy=sockets.target
