# sys-apps/my-systemd-units/files/services-basic/radeon-power_profile_at.service
[Unit]
Description=Set Radeon Power Profile to %I

[Service]
Type=oneshot
ExecStart=/usr/sbin/radeon-power-profile %I

[Install]
WantedBy=multi-user.target
