# net-wireless/a2dp-switch/files/80-bt-headset-a2dp.rules
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/a2dp-fix-wrapper"
